# deploy/linode/systemd/fleximos-foundation.service
[Unit]
Description=Fleximotion Foundation API (Identity + Amoeba)
After=network.target

[Service]
Type=simple
User=fleximos
WorkingDirectory=/srv/fleximos
EnvironmentFile=/etc/fleximos/fleximos.env
Environment=PORT=4010
Environment=HOST=127.0.0.1
Environment=FLEXI_DB_DIR=/var/lib/fleximos/foundation-pglite
ExecStart=/usr/bin/node apps/api-foundation/server.mjs
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
