# deploy/linode/systemd/fleximos-foundation.service
# systemd USER unit — installed to ~/.config/systemd/user by install.sh,
# which replaces __NODE_BIN__ with the absolute path of your node binary.
[Unit]
Description=Fleximotion Foundation API (Identity + Amoeba)
After=network.target

[Service]
Type=simple
WorkingDirectory=%h/fleximos
EnvironmentFile=%h/fleximos-data/fleximos.env
Environment=PORT=4010
Environment=HOST=127.0.0.1
Environment=FLEXI_DB_DIR=%h/fleximos-data/foundation-pglite
ExecStart=__NODE_BIN__ apps/api-foundation/server.mjs
Restart=always
RestartSec=3

[Install]
WantedBy=default.target
