# deploy/linode/systemd/fleximos-ops-scheduler.timer
[Unit]
Description=Run the Fleximotion Ops scheduler every minute

[Timer]
OnCalendar=*-*-* *:*:00
AccuracySec=5s
Persistent=false

[Install]
WantedBy=timers.target
